# tea-testing-baseline-nvidia/airootfs/etc/xdg/autostart/tealinux-installer.desktop
[Desktop Entry]
Categories=Development;
Exec=tealinux-installer
Icon=tealinux-installer
Name=tealinux-installer
Terminal=false
Type=Application
